# files/hbase-file-increase.service
[Unit]
Description=HBase increase the limit

[Service]
Type=simple
User=root

ExecStart=/usr/bin/hbase-file-increase

[Install]
WantedBy=multi-user.target
