# files/hbase-hdfs.service
[Unit]
Description=HBase HDFS

[Service]
Type=forking
User=root

ExecStart=/usr/bin/hbase-hdfs

[Install]
WantedBy=multi-user.target
